# flathunter-callback-handler.service
[Unit]
Description=Flathunter Telegram Callback Handler
After=network.target

[Service]
Type=simple
User=user
WorkingDirectory=/home/user/Flatjaga
ExecStart=/usr/bin/python3 /home/user/Flatjaga/scripts/run_callback_handler.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment
Environment="PYTHONUNBUFFERED=1"

[Install]
WantedBy=multi-user.target
